# ansible/roles/tg_bot/templates/tg_bot.service.j2
[Unit]
Description=Telegram Bot Service
After=network.target

[Service]
User=root
ExecStart=/opt/tg-bot/venv/bin/python3 /opt/tg-bot/bot.py
Restart=always

[Install]
WantedBy=multi-user.target
